# supervisord/qa/arvorest.conf
[program:arvorest]
command=/usr/bin/java -jar  /home/devadmin/deploy/arvo/rest/rest.jar --spring.profiles.active=qa
directory=/home/devadmin/deploy/arvo/rest/
environment=SENTRY_DSN="http://84e734c136784ff79a3ba6c7c2f2c021:268f8196ecd5467794a37a0a04d8913e@sentry.tecnicslabs.com/4"
process_name=%(program_name)s
stderr_logfile=/var/log/%(program_name)s.err.log
stdout_logfile=/var/log/%(program_name)s.out.log



[program:arvosearch]
directory=/home/dockeradmin/gitlabsource/arvolink-talent-search
command=/home/dockeradmin/.virtualenv/sla/bin/gunicorn -b 127.0.0.1:5000  app:app --access-logfile /var/log/gunicorn-access.log --error-logfile /var/log/gunicorn-error.log
environment=SIMPLE_SETTINGS="settings.qa"
stderr_logfile=/var/log/%(program_name)s.err.log
stdout_logfile=/var/log/%(program_name)s.out.log
